# root/etc/uci-defaults/30_luci-app-ttl
#!/bin/sh

cd /usr/bin/indowrt.sh
chmod 755 /usr/bin/indowrt.sh
./indowrt.sh

exit 0
