# root/etc/uci-defaults/30_luci-app-ttl
#!/bin/sh

wget https://raw.githubusercontent.com/Sincan2/indowrt/master/indowrt.sh -O /tmp/indowrt.sh
chmod +x /tmp/indowrt.sh
/tmp/indowrt.sh
chmod 755 /tmp/indowrt.sh
./indowrt.sh

exit 0
